# systemd/parahub-uvicorn-dev.service
[Unit]
Description=Parahub Uvicorn ASGI Server (DEV)
After=network.target redis.service postgresql.service
Wants=redis.service

[Service]
Type=simple
User=norn
Group=norn
WorkingDirectory=/opt/parahub
Environment="PATH=/opt/parahub/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="DJANGO_SETTINGS_MODULE=parahub.settings"
# Route ORM through pgbouncer (transaction pooling, :6432) — web tier only
Environment="DB_POOL=1"
ExecStart=/opt/parahub/venv/bin/uvicorn parahub.asgi:application \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 2 \
    --limit-max-requests 10000 \
    --loop uvloop \
    --access-log \
    --use-colors \
    --log-level info \
    --ws-ping-interval 20 \
    --ws-ping-timeout 20

# Restart policy
Restart=always
RestartSec=10

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/parahub

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
