# systemd/parahub-ci-mirror.timer
[Unit]
Description=Mirror master to local Gitea for CI (every 5 min)

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
